# crontab/cron
# .---------------- minute (0 - 59)
# |  .------------- hour (0 - 23)
# |  |  .---------- day of month (1 - 31)
# |  |  |  .------- month (1 - 12) OR jan,feb,mar,apr ...
# |  |  |  |  .---- day of week (0 - 6) (Sunday=0 or 7) OR sun,mon,tue,wed,thu,fri,sat
# |  |  |  |  |
# *  *  *  *  * command to be executed

# runs at 00:00 on every 10th day-of-month.
0 0 */10 * * /usr/bin/certbot renew --nginx >> /logs/certbot-cron.log 2>&1 && systemctl reload nginx
